# hbakd.initd
#!/sbin/openrc-run

name="${RC_SVCNAME}"

command="/usr/bin/hbakd"
command_args=""

depend() {
	need net
	provide hbakd
}
